# update linux server.service
[Unit]
Description=Update linux servers

[Service]
ExecStart=Path to bash-script to update servers