# update linux server.timer
[Timer]
Unit=kiiolabs.service
OnCalendar=*-*-14

[Install]
WantedBy=timers.target